# packages/vm-agent/agent-service
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Takaro agent"
supervisor="supervise-daemon"
command="/bin/agent"
pidfile="/run/agent.pid"

depend() {
	after net
}
